# deployments/systemd/djangosite.socket
[Unit]
Description=gunicorn socket (djangosite)

[Socket]
ListenStream=/run/gunicorn/djangosite.sock

[Install]
WantedBy=sockets.target